# roles/vhosts/acp_server_codex/templates/acp-bridge.service.j2
[Unit]
Description=XWorkmate Codex ACP bridge server
After=network-online.target {{ acp_codex_service_name }}.service
Wants=network-online.target

[Service]
Type=simple
User={{ acp_codex_service_user }}
Group={{ acp_codex_service_group }}
WorkingDirectory={{ acp_codex_workdir }}
Environment=HOME={{ acp_codex_workdir }}
Environment=TERM=xterm-256color
Environment=ACP_LISTEN_ADDR={{ acp_codex_bridge_listen_host }}:{{ acp_codex_bridge_listen_port }}
Environment=ACP_ALLOWED_ORIGINS={{ acp_codex_bridge_allowed_origins | join(',') }}
{% if acp_codex_auth_token | trim | length > 0 %}
Environment=ACP_AUTH_TOKEN={{ acp_codex_auth_token }}
{% endif %}
{% for key, value in acp_codex_environment | dictsort %}
{% if value | string | trim | length > 0 %}
Environment={{ key }}={{ value }}
{% endif %}
{% endfor %}
ExecStart={{ acp_codex_bridge_binary_path }} serve --listen {{ acp_codex_bridge_listen_host }}:{{ acp_codex_bridge_listen_port }}
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
