# roles/vhosts/deploy_acp_vhosts/templates/xworkmate-bridge.service.j2
[Unit]
Description=XWorkmate Bridge ACP Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ xworkmate_bridge_service_user }}
Group={{ xworkmate_bridge_service_group }}
WorkingDirectory={{ xworkmate_bridge_workdir }}
ExecStart={{ xworkmate_bridge_binary_path }} serve --listen {{ xworkmate_bridge_listen_host }}:{{ xworkmate_bridge_listen_port }}
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
