# roles/vhosts/acp_bridge_server/templates/acp-bridge-server.service.j2
[Unit]
Description=XWorkmate ACP Bridge Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ acp_bridge_server_service_user }}
Group={{ acp_bridge_server_service_group }}
WorkingDirectory={{ acp_bridge_server_workdir }}
ExecStart={{ acp_bridge_server_binary_path }} serve --config /etc/xworkmate/xworkmate-acp-bridge-server.yaml --listen {{ acp_bridge_server_listen_host }}:{{ acp_bridge_server_listen_port }}
Restart=on-failure
RestartSec=2
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
