# roles/vhosts/acp_server_hermes/templates/hermes-acp-adapter.service.j2
[Unit]
Description=XWorkmate Hermes ACP adapter
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ acp_hermes_service_user }}
Group={{ acp_hermes_service_group }}
WorkingDirectory={{ acp_hermes_workdir }}
Environment=HOME={{ acp_hermes_home }}
Environment=TERM=xterm-256color
Environment=HERMES_ADAPTER_LISTEN_ADDR={{ acp_hermes_listen_host }}:{{ acp_hermes_listen_port }}
Environment=HERMES_ADAPTER_BIN={{ acp_hermes_binary_path }}
Environment=HERMES_ADAPTER_ARGS={{ acp_hermes_args }}
Environment=HERMES_ADAPTER_ALLOWED_ORIGINS={{ acp_hermes_allowed_origins | join(',') }}
{% for key, value in acp_hermes_environment | dictsort %}
Environment={{ key }}={{ value }}
{% endfor %}
ExecStart={{ acp_hermes_bridge_binary_path }} hermes-acp-adapter --listen {{ acp_hermes_listen_host }}:{{ acp_hermes_listen_port }} --hermes-bin {{ acp_hermes_binary_path }} --hermes-args "{{ acp_hermes_args }}"
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
